# supervisor-app-staging.conf
[supervisord]
loglevel = INFO
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[unix_http_server]
file = /project/django_gunicorn.sock
chmod = 0700
username = dummy
password = dummy

[program:gunicorn]
command=gunicorn django_gunicorn.wsgi.%(ENV_SERVICE_ENV)s -c /project/gunicorn_conf.py
directory=/project/django_gunicorn
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true


[program:nginx]
command = /usr/sbin/nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

